# typofixer.service
[Unit]
Description=Streamlit app to correct typos with an LLM.
After=network.target

[Service]
WorkingDirectory=/srv/typofixer
ExecStart=make run

[Install]
WantedBy=multi-user.target
